# scripts/launchd/com.senna.postmortem-2026-05-01.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.senna.postmortem-2026-05-01</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/macbookair/Documents/VibeCoding/senna-brands/sennabrands-realtime-monitoring/scripts/postmortem-2026-05-01.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/macbookair/Documents/VibeCoding/senna-brands/sennabrands-realtime-monitoring</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Month</key>
    <integer>5</integer>
    <key>Day</key>
    <integer>2</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>7</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/senna-postmortem-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/senna-postmortem-stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
